# services/websocket2.service
[Unit]
Description=Websocket 2 service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=hs
ExecStart=/home/hs/scripts/websocket2.sh

[Install]
WantedBy=multi-user.target
